# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ashpy]
command=/opt/ashpy/venv/bin/ashpy serve --bind 127.0.0.1:50051
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:ash-serve]
# M4: `ash serve` now hosts the Rust QueryHost gRPC server on :50052
# that the Python FastAPI layer calls from its /v1/chat handler.
command=/usr/local/bin/ash serve --host 127.0.0.1 --port 50052 --sidecar http://127.0.0.1:50051
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
